# imageroot/systemd/user/dhgen.service
[Unit]
Description=Diffie-Hellman group generator
Requisite=ejabberd.service

[Service]
Type=oneshot
ExecStart=runagent generate-dhpem
ExecStartPost=runagent install-dh-certificate-key
WorkingDirectory=%S/state
SyslogIdentifier=%u/%N
