# imageroot/systemd/user/purge-mnesia-database.timer
[Unit]
Description=Timer of purge mnesia database

[Timer]
OnActiveSec=15 minutes
OnUnitInactiveSec=1 days

[Install]
WantedBy=timers.target
